# deploy/agent.service
[Unit]
Description=MySQL PITR Agent
Documentation=https://github.com/a-shan/mysql-pitr
After=network.target
Wants=mysql.service

[Service]
Type=simple
ExecStart=/usr/local/bin/mysql-pitr-agent flashback --config=/etc/agent/config.json
Restart=always
RestartSec=5
EnvironmentFile=-/etc/agent/env
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
